# service/calyx_daemon.service
[Unit]
Description=Calyx Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/calyx
ExecStart=/usr/bin/python3 /opt/calyx/tools/calyx_daemon.py
Restart=on-failure
User=calyx

[Install]
WantedBy=multi-user.target
